# deploy/holo-rick-auto-update.timer
[Unit]
Description=Run Holo Rick auto update periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=3min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
